# etc/hdfs-auto-snapshot.cron.monthly
#!/bin/sh
su hdfs -c 'zfs-auto-snapshot --quiet --syslog --label=monthly --keep=12 //'
