# etc/hdfs-auto-snapshot.cron.hourly
#!/bin/sh
su hdfs -c 'hdfs-auto-snapshot --quiet --syslog --label=hourly --keep=24 //'
