# etc/hdfs-auto-snapshot.cron.weekly
#!/bin/sh
su hdfs -c 'hdfs-auto-snapshot --quiet --syslog --label=weekly --keep=8 //'
