# etc/hdfs-auto-snapshot.cron.daily
#!/bin/sh
su hdfs -c 'hdfs-auto-snapshot --quiet --syslog --label=daily --keep=31 //'
